# omablue/config/systemd/user/omasync@.timer
[Unit]
Description=Omasync Timer — %i

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
